# django/provision/roles/5_gunicorn/templates/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ GUNICORN_BIND }}

[Install]
WantedBy=sockets.target